# rc.d/backer
#!/bin/sh

# REQUIRE: DAEMON
# PROVIDE: backer
# KEYWORD: shutdown

export PATH=/usr/local/bin:$PATH
export PYTHONPATH=/usr/local/backer/lib

. /etc/rc.subr

name="backer"
load_rc_config $name

rcvar="${name}_enable"
pidfile="/var/run/$name.pid"
command="/usr/sbin/daemon"
command_args="-c -r -t $name -S -T $name -P $pidfile /usr/local/backer/bin/backer"

run_rc_command "$1"
